# releng/airootfs/etc/systemd/system/getty@tty2.service.d/autologin.conf
# Root autologin console on tty2 — fallback for debugging the live ISO.
# tty1 is used by GDM/Xorg; use Ctrl+Alt+F2 to reach this console.
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --noclear %I 38400 linux
