# releng/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
# tty1 is reserved for GDM/Xorg on the live ISO.
# Disable getty on tty1 so GDM can freely acquire VT1 for the Xorg session.
# A root autologin fallback console is available on tty2.
[Unit]
ConditionPathExists=!/etc/gdm/custom.conf
