# releng/airootfs/etc/systemd/system/gdm.service.d/autologin.conf
# Garante que /run/nologin já foi removido antes do GDM iniciar.
#
# O GDM verifica /run/nologin independentemente do PAM antes de tentar
# o autologin. Se o GDM arrancar antes do systemd-user-sessions.service
# remover esse ficheiro, o autologin falha silenciosamente e o GDM
# apresenta o ecrã de login manual.
#
# Num live ISO, o GDM tende a ganhar a corrida contra esse serviço.
# Este drop-in resolve o problema declarando a dependência explícita.
[Unit]
After=systemd-user-sessions.service
Wants=systemd-user-sessions.service
